# ~/Library/LaunchAgents/com.user.tennisbooking.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.tennisbooking</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/kohanchen/Documents/projects/tennis/run_tennis_booking.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>11</integer>
        <key>Minute</key>
        <integer>57</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/kohanchen/Documents/projects/tennis/scheduler.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kohanchen/Documents/projects/tennis/scheduler_error.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist> 